# deploy/samizdat.service.tmpl
[Unit]
Description=Samizdat server (@INSTANCE@)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Environment=HOME=@HOME@
ExecStart=@BIN@ serve --config @CFG@ --webdir @WEBDIR@@EXTZIP@
WorkingDirectory=@REPO@
Restart=on-failure
RestartSec=3
# Bind :443 without running as root (needed for TLS post-M1):
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
