# deploy/samizdat.service.tmpl
[Unit]
Description=Samizdat server (@INSTANCE@)
After=network-online.target

[Service]
Type=simple
# Secrets (ANTHROPIC_API_KEY, OPENAI_API_KEY, …) — user services do NOT source
# ~/.bashrc, so the shell env is invisible here. Load them from a gitignored
# .env in the repo root instead. The leading `-` makes the file optional.
EnvironmentFile=-@REPO@/.env
ExecStart=@BIN@ serve --config @CFG@ --webdir @WEBDIR@@EXTZIP@
WorkingDirectory=@REPO@
Restart=on-failure
RestartSec=3
# User service (systemctl --user) so restarts need no sudo. Lingering keeps it
# running across logout/reboot (loginctl enable-linger). Runs as the login user.
#
# TLS/:443 is post-M1; the server currently binds an unprivileged port over
# plain HTTP. When TLS lands, a user service cannot hold CAP_NET_BIND_SERVICE —
# grant low-port binding once with either:
#   sudo sysctl -w net.ipv4.ip_unprivileged_port_start=443   (+ /etc/sysctl.d)
#   or sudo setcap cap_net_bind_service=+ep @BIN@             (re-run after each build)

[Install]
WantedBy=default.target
